# keepPrimaryDMZRunning.service
[Unit]
Description=Service that keeps running the echo-server from startup.
After=network.target
StartLimitIntervalSec=0
[Install]
WantedBy=multi-user.target
[Service]
Type=simple
ExecStart=/home/kevin/Desktop/Git/IT-490/primaryDMZRabbitMQServer.php
Restart=always
RestartSec=5
